# vanurses/vanurses-api.service
[Unit]
Description=VANurses.com API Service
After=network.target postgresql.service

[Service]
Type=simple
User=ian
WorkingDirectory=/home/ian/vanurses/api
Environment=PATH=/home/ian/vanurses/api/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/ian/vanurses/api/venv/bin/gunicorn app.main:app --workers 3 --worker-class uvicorn.workers.UvicornWorker --bind 0.0.0.0:5011 --timeout 120
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
